# scripts/com.tejasmahadik.asr-webui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tejasmahadik.asr-webui</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tejasmahadik/Documents/agentic-lecture-notes/venv/bin/python</string>
        <string>/Users/tejasmahadik/Documents/agentic-lecture-notes/web_ui/app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tejasmahadik/Documents/agentic-lecture-notes</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/tejasmahadik/Library/Logs/asr_webui_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tejasmahadik/Library/Logs/asr_webui_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
